# overlay/etc/systemd/system/sshd.service.d/foldingos.conf
[Unit]
Requires=foldingos-ssh-provision.service
After=foldingos-ssh-provision.service
ConditionPathExists=/data/config/ssh/authorized_keys

[Service]
ExecStartPre=
