# overlay/usr/lib/systemd/system/var-log-journal.mount
[Unit]
Description=FoldingOS persistent journal
Requires=foldingos-journal-dir.service
After=foldingos-journal-dir.service
Before=foldingos-journal-flush.service

[Mount]
What=/data/logs/journal
Where=/var/log/journal
Type=none
Options=bind

[Install]
WantedBy=local-fs.target
