# overlay/etc/systemd/system/systemd-networkd-wait-online.service.d/foldingos.conf
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any --ipv4 --ignore=lo --timeout=120
TimeoutStartSec=120
